# misc/systemd/ioc.service
[Unit]
Description=IOC

[Service]
ExecStart=/usr/bin/sh /opt/run_ioc.sh
Restart=always

[Install]
WantedBy=multi-user.target
